# scripts/ops/launchd/com.assyntrax.daily-publish.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.assyntrax.daily-publish</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>scripts/ops/run_daily_publish_site.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__REPO_PATH__</string>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>19</integer>
    <key>Minute</key>
    <integer>10</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>__REPO_PATH__/results/ops/logs/daily_publish.out.log</string>
  <key>StandardErrorPath</key>
  <string>__REPO_PATH__/results/ops/logs/daily_publish.err.log</string>

  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>

